# root/etc/services.d/xvfb/run
#!/usr/bin/execlineb -P

/usr/bin/Xvfb :1