# root/etc/services.d/xvfb/run
#!/usr/bin/execlineb -P

export DISPLAY=:1 Xvfb :1